# try-daemon.service
[Unit]
Description=Foo

[Service]
ExecStart=/usr/sbin/try-daemon

[Install]
WantedBy=multi-user.target


